# hosts/draco-router/rootfs/etc/cron.d/registryUpdater
*/1 * * * *	root	/usr/bin/registryUpdater /etc/registryUpdater/*.conf
