# hosts/draco-router/rootfs/etc/rc.local
#!/bin/bash

mount /etc/network
mount /etc/openvpn
mount /etc/shorewall
mount /etc/registryUpdater

/etc/init.d/shorewall restart
/etc/init.d/networking restart
/etc/init.d/openvpn restart
